# daemon/conflip@.service.in
[Unit]
Description=Conflip Daemon

[Service]
Type=simple
ExecStart=%{INSTALL_BINS}/conflipd --systemd-log --slice %I
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-abnormal

[Install]
WantedBy=default.target
