# daemon/conflip.service.in
[Unit]
Description=Conflip Service
Documentation=https://github.com/Skycoder42/Conflip

[Service]
Type=notify
NotifyAccess=exec
ExecStart=%{INSTALL_BINS}/conflipd -platform minimal --backend systemd %{SLICE}
ExecReload=%{INSTALL_BINS}/conflipd -platform minimal --backend systemd reload $MAINPID %{SLICE}
ExecStop=%{INSTALL_BINS}/conflipd -platform minimal --backend systemd stop $MAINPID %{SLICE}
#WatchdogSec=10
Restart=on-abnormal
RuntimeDirectory=%{conflip}

[Install]
WantedBy=default.target
